# files/etc/profile.d/tmux.sh
alias tmux="/usr/bin/env tmux -u -2 -f `[ -f ~/.tmux.conf ] && echo "~/.tmux.conf" || echo "/etc/tmux.conf"`"
